# deployment/systemd/llm-gateway.service
[Unit]
Description=Claude Universal Custom Proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=llm-gateway
Group=llm-gateway
WorkingDirectory=/opt/claude-universal-custom-proxy
EnvironmentFile=/etc/claude-universal-custom-proxy/gateway.env
ExecStart=/opt/claude-universal-custom-proxy/.venv/bin/uvicorn llm_proxy_gateway.main:app --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/claude-universal-custom-proxy

[Install]
WantedBy=multi-user.target
