# deployment/systemd/llm-gateway.service
[Unit]
Description=OpenAI Compatible LLM Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=llm-gateway
Group=llm-gateway
WorkingDirectory=/opt/openai-compatible-llm-gateway
EnvironmentFile=/etc/openai-compatible-llm-gateway/gateway.env
ExecStart=/opt/openai-compatible-llm-gateway/.venv/bin/uvicorn llm_proxy_gateway.main:app --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/openai-compatible-llm-gateway

[Install]
WantedBy=multi-user.target
